# etc/crontabs/root
0 2 * * * /usr/share/shadowsocksr/ssrplusupdate.sh
